# etc/udev/rules.d/99-k380.rules
KERNEL=="hidraw*", ACTION=="add", SUBSYSTEM=="hidraw",  RUN+="/sbin/k380_conf -d /dev/%k -f on"

